# docker/shadowban-alerts/crontab
*/1 * * * * cd /app; python ShadowbanAlerts.py